# service/financetracker.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Task Scheduler task definition for the FinanceTracker supervisor.

  One task keeps BOTH always-on servers alive: it runs service/supervisor.py
  under pythonw.exe, which relaunches the FastAPI backend (8010) and the static
  PWA server (4173) whenever either stops listening.

  The register-task.ps1 script replaces the __REPO_ROOT__ placeholder with the
  absolute path to the repository before registering this task.  Do NOT edit the
  placeholder manually; run the register script instead.

  Windowless: pythonw.exe is a GUI-subsystem binary, so no console window ever
  exists — unlike powershell.exe -WindowStyle Hidden, which Windows Terminal
  (the Windows 11 default console host) can still show.

  Security: RunLevel=LeastPrivilege — the task runs as the logged-in user without
  elevation, which means the servers keep running while the laptop is locked
  (InteractiveToken user task).  No UAC prompt, no admin rights.

  No secrets appear in this file.  All sensitive config is loaded from .env at
  runtime by python-dotenv inside the Python process.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>FinanceTracker — windowless supervisor for the backend and PWA servers (starts on login and on unlock, self-heals crashes)</Description>
    <Version>2.0.0</Version>
  </RegistrationInfo>

  <Triggers>
    <!-- Start automatically when the user logs in (FR-6). -->
    <LogonTrigger>
      <Enabled>true</Enabled>
    </LogonTrigger>
    <!-- Also fire on session unlock: if the supervisor died while the machine
         was locked/asleep, unlocking revives it.  When it is already running,
         MultipleInstancesPolicy=IgnoreNew makes this a no-op. -->
    <SessionStateChangeTrigger>
      <Enabled>true</Enabled>
      <StateChange>SessionUnlock</StateChange>
    </SessionStateChangeTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <!-- Run as the interactively logged-in user (keeps running when screen is locked). -->
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <!-- Only one instance allowed; extra triggers are ignored. -->
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <!-- Battery state must not stop the service. -->
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <!-- Run indefinitely — no time limit. -->
    <ExecutionTimeLimit>PT0S</ExecutionTimeLimit>
    <!-- Start immediately if the trigger was missed (e.g. machine was off). -->
    <StartWhenAvailable>true</StartWhenAvailable>
    <!-- Auto-restart on crash: up to 3 attempts, 1 minute between each. -->
    <RestartOnFailure>
      <Interval>PT1M</Interval>
      <Count>3</Count>
    </RestartOnFailure>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <Command>__REPO_ROOT__\venv\Scripts\pythonw.exe</Command>
      <Arguments>"__REPO_ROOT__\service\supervisor.py"</Arguments>
      <WorkingDirectory>__REPO_ROOT__</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
